# debian/api-bitbucket/DEBIAN/postinst
#!/bin/bash
systemctl enable api-bitbucket.service
systemctl start api-bitbucket.service
