# debian/api-bitbucket/DEBIAN/preinst
#!/bin/bash
touch /tmp/foo
systemctl stop api-bitbucket.service || true
systemctl disable api-bitbucket.service || true
